# deploy_tools/gunicorn-upstart.template.conf
description "Gunicorn server for SITENAME"

start on net-device-up
stop on shutdown

respawn

setuid pando
chdir /home/pando/sites/SITENAME/source

exec gunicorn \
    --bind unix:/tmp/SITENAME.socket \
    superlists.wsgi:application
